# ansible-autofix/roles/autofix/templates/cpu_autoheal.service.j2
[Unit]
Description=Run Cpu_Autoheal Script

[Service]
Type=oneshot
ExecStart=/usr/local/sre-autofix/cpu_autoheal.sh
