# ansible-autofix/roles/autofix/templates/deadlock-autoheal.service.j2
[Unit]
Description=Run Deadlock Auto-Heal Script

[Service]
Type=oneshot
ExecStart=/usr/local/sre-autofix/deadlock_autoheal.sh

