# ansible-autofix/roles/autofix/templates/disk-autoheal.service.j2
[Unit]
Description=Run Disk Auto-Heal Script

[Service]
Type=oneshot
ExecStart=/usr/local/sre-autofix/disk_autoheal.sh

