# ansible-autofix/roles/autofix/templates/memory_autoheal.service.j2
[Unit]
Description=Run Memory_Autoheal Script

[Service]
Type=oneshot
ExecStart=/usr/local/sre-autofix/memory_autoheal.sh
