# ansible-autofix/roles/autofix/templates/network_autoheal.service.j2
[Unit]
Description=Run Network_Autoheal Script

[Service]
Type=oneshot
ExecStart=/usr/local/sre-autofix/network_autoheal.sh
