# ansible-autofix/roles/autofix/templates/system_autoheal.service.j2
[Unit]
Description=Run System_Autoheal Script

[Service]
Type=oneshot
ExecStart=/usr/local/sre-autofix/system_autoheal.sh
